# res/blocker.conf
start on (started docker)

respawn
console log

env AWS_SHARED_CREDENTIALS_FILE=/etc/blocker/.aws/credentials

exec /usr/local/bin/blocker

