# conf/blocker.conf
start on (started docker)

respawn
console log

exec /usr/local/bin/blocker

